# ninjacore/deploy/ninjacore.service
[Unit]
Description=ninjacore — NinjaDispute API (Axum + Tokio)
After=network.target surrealdb.service
Wants=surrealdb.service

[Service]
Type=simple
User=ninjacore
Group=ninjacore
WorkingDirectory=/opt/ninjacore
EnvironmentFile=/opt/ninjacore/.env
# TCP for in-host browser-runner callbacks; Unix socket for Caddy.
Environment=BIND_ADDR=127.0.0.1:8080
Environment=UNIX_SOCKET=/run/ninjacore/ninjacore.sock
# 0o660 + caddy in the ninjacore group is the locked-down setup.
Environment=UNIX_SOCKET_MODE=0o660
# Pre-create the socket directory with the right group ownership.
RuntimeDirectory=ninjacore
RuntimeDirectoryMode=0750
# If caddy runs as user `caddy`, make sure `caddy` is in group `ninjacore`:
#   usermod -a -G ninjacore caddy
ExecStart=/opt/ninjacore/bin/ninjacore
Restart=on-failure
RestartSec=2
LimitNOFILE=65536
# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/ninjacore
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target
