# apps/ninjacore-web-svelte/ninjacore-web.service
[Unit]
Description=NinjaCore SvelteKit Frontend (Bun)
After=network.target
Wants=ninjacore-api.service

[Service]
Type=simple
User=root
WorkingDirectory=/home/user/ninjacore/apps/ninjacore-web-svelte

# Bun runtime
ExecStart=/root/.bun/bin/bun run start

# Performance: Pin to CPU cores 0-3 (Axum uses 4-7)
CPUAffinity=0-3

# Resource limits
MemoryLimit=512M
MemoryAccounting=yes

# Auto-restart on failure
Restart=on-failure
RestartSec=5s

# Logging
StandardOutput=journal
StandardError=journal

# Environment
Environment="NODE_ENV=production"
Environment="RUST_API_ORIGIN=unix:///tmp/ninjacore.sock"
Environment="API_PORT=3100"

[Install]
WantedBy=multi-user.target
